# ops/supervisord.conf
; %(ENV_VARS)s is an operator-supplied placeholder for environment variables.
[supervisord]
nodaemon=true
logfile=/var/log/aindy/supervisord.log
pidfile=/tmp/supervisord.pid

[program:aindy-web]
command=uvicorn AINDY.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=%(ENV_VARS)s

[program:aindy-worker]
command=python -m AINDY.worker.worker_loop
directory=/app
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
environment=WORKER_HEALTH_PORT=8001,%(ENV_VARS)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:aindy-memory-ingest]
command=python -m AINDY.worker.memory_ingest_worker
directory=/app
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
environment=WORKER_HEALTH_PORT=8002,%(ENV_VARS)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:aindy-metric-writer]
command=python -m AINDY.worker.metric_writer_worker
directory=/app
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=15
environment=WORKER_HEALTH_PORT=8003,%(ENV_VARS)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[group:aindy]
programs=aindy-web,aindy-worker,aindy-memory-ingest,aindy-metric-writer
